# service/install/macos/com.knws.screencontrol.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.knws.screencontrol.service</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/ScreenControlService</string>
        <string>--port</string>
        <string>3459</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Library/Logs/ScreenControl/service-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Logs/ScreenControl/service-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/bin</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
